# servers/portal/crontab
0 4 * * * kirios server portal-backup-home | systemd-cat -t backup-home
5 4 * * 0 kirios server portal-backup-portainer | systemd-cat -t backup-portainer