# servers/portal/crontab
0 4 * * * kirios server portal-backup-vikunja
5 4 * * 0 kirios server portal-backup-portainer
10 4 * * 0 kirios server portal-backup-portainer